# gui/.config/systemd/user/dropbox.service
[Unit]
Description=Dropbox
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/dropbox

[Install]
WantedBy=sway-session.target
